# .ddev/web-build/espocrm.cron
* * * * * www-data cd /var/www/html && php cron.php > /dev/null 2>&1
